# roles/extractor/scripts/launchd/com.extractor.session-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.extractor.session-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexander/Github/FMT-exocortex-template/roles/extractor/scripts/session-watcher.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/alexander/logs/extractor/session-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alexander/logs/extractor/session-watcher-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/alexander</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
